# kubernetes_masta/etcd.service
[Unit]
Description=Etcd Server
After=network.target
 
[Service]
Type=simple
StandardOutput=null
WorkingDirectory=/opt/etcd
ExecStart=/opt/etcd/etcd
 
[Install]
WantedBy=multi-user.target